# big-stop-akonadi/etc/xdg/autostart/big_stop_akonadi.desktop
[Desktop Entry]
Name=finaliza o akonadi se nao estiver em uso
Comment=
Icon=
Exec=big_stop_akonadi
Terminal=false
Type=Application
OnlyShowIn=KDE;