# bigcontrolcenter-system-optimizes/etc/xdg/autostart/big-nm-applet-update-notifier.desktop
[Desktop Entry]
Name=Update Notifier
Comment=Update Notifier
Icon=update-notifier
Exec=bigupdate-notifier&
Terminal=false
Type=Application