# bigdesktop-icons/etc/xdg/autostart/big-nm-applet-desktopicons.desktop
[Desktop Entry]
Name=Desktop Icons
Comment=Trash Computer and Network
Icon=computer
Exec=bigdesktop-icons
Terminal=false
Type=Application