# bigcontrolcenter-system-optimizes/etc/init/bigswap.conf
# BigBoot

description	"Tweaks for BigLinux"
author		"Bruno Goncalves Araujo"

start on (filesystem)
script

	/usr/bin/bigswap

end script
