# transparent-hugepage-debug-disable/etc/init/transparent_hugepage_debug_disable.conf
# hugepage debug disable

description	"hugepage debug disable"
author		"Bruno Goncalves Araujo"

start on (filesystem)
script
	echo "never" > /sys/kernel/mm/transparent_hugepage/defrag
end script
