# broffice.org-big/etc/profile.d/libreoffice.sh
export OOO_FORCE_DESKTOP="$(cat ~/.oooforce)"