# bigsudo/etc/xdg/autostart/X11forRoot.desktop
[Desktop Entry]
Name=X11 for Root
Comment=X11 for Root
Exec=xhost +SI:localuser:root
Terminal=false
Type=Application