# bigcontrolcenter-hardware-keyboard/etc/xdg/autostart/big-nm-applet-touchfreeze.desktop
[Desktop Entry]
Name=Touchfreeze
Comment=Touchpad freeze
Exec=bigtouchfreeze&
Terminal=false
Type=Application
Icon=touchpad