# bigcontrolcenter-hardware-keyboard/etc/xdg/autostart/big-nm-applet-kvkbd.desktop
[Desktop Entry]
Name=Virtual Keyboard
Comment=Virtual Keyboard
Exec=bigkvkbd-applet&
Terminal=false
Type=Application
Icon=kvkbd