# bigcontrolcenter-hardware-keyboard/etc/xdg/autostart/big-nm-applet-gsynaptics.desktop
[Desktop Entry]
Name=Touchpad
Comment=Set your touchpad preferences
Exec=gsynaptics-init --sm-disable&
Terminal=false
Type=Application
Icon=touchpad