# bigcontrolcenter-hardware-keyboard/etc/xdg/autostart/big-nm-applet-eviacam.desktop
[Desktop Entry]
Name=eViacam
Comment=eViacam
Exec=bigeviacam-applet&
Terminal=false
Type=Application
Icon=eviacam