# web-quest/etc/xdg/autostart/web-quest.desktop
[Desktop Entry]
Exec=web-quest
Name=Web-Quest
Type=Application
